# deploy/s6/stoat/finish
#!/bin/bash
cd /srv/stoat
podman compose down

