# deploy/s6/stoat/run
#!/bin/bash
set -e
cd /srv/stoat
exec podman compose up 2>&1

