# deploy/s6/stoat-admin/finish
#!/bin/bash
cd /srv/stoat-admin
podman compose down

